# units/webkit-search-indexer-run@.service
[Unit]
Description=WebKitSearch Indexer Run (%i)
After=webkit-search-indexer-setup@%i.service
Requires=webkit-search-indexer-setup@%i.service

[Service]
Type=oneshot
User=%i
Group=%i
WorkingDirectory=/home/%i
RuntimeDirectory=%i
Environment=KEEP_WORKING=1
EnvironmentFile=/home/%i/environment

ExecStart=/home/%i/${MOZSEARCH}/infrastructure/indexer-run.sh ${CONFIG} {INDEX}
